# docker/rootfs/etc/s6-overlay/s6-rc.d/init-usermod/run
#!/usr/bin/execlineb -P

with-contenv

UID=${UID:-1001}
GID=${GID:-1001}

s6-echo "User Id $UID, Group id $GID"
usermod -u "$UID" nodejs && groupmod -g "$GID" nodejs
