# docker/rootfs/etc/s6-overlay/s6-rc.d/init-usermod/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

PUID=${PUID:-1001}
PGID=${PGID:-1001}

echo "UID=${PUID}, GID=${PGID}"

usermod -u "${PUID}" nodejs && groupmod -g "${PGID}" nodejs
