# docker/rootfs/etc/s6-overlay/s6-rc.d/dashboard/finish
#!/usr/bin/execlineb -P

s6-echo "Stopping dashboard"
